# common/projects/system_service_examples/.service
[Unit]
Description=Sysprog service
Wants=network-online.target
After=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
Type=simple
ExecStart=/home/ubuntu/ML/anaconda3/bin/python3.7 /YOUR_ABSOLUTE_PATH/sysprog_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target